# examples/data/detailed.build_context/usr/local/lib/systemd/system/p5-systemd4docker-example.service
[Unit]
Description = p5 systemd4docker example service

[Service]
Type = oneshot
ExecStart = sleep 42
RemainAfterExit = yes

[Install]
WantedBy = multi-user.target
